# docker/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/with-contenv sh
# shellcheck shell=sh

exec busybox crond -f -l 0 -L /dev/stdout
