# docker/var/spool/cron/crontabs/seedboxsync
# min   hour    day     month   weekday command
*       *       *       *       *       seedboxsync -q sync blackhole --ping
*/15    *       *       *       *       seedboxsync -q sync seedbox --ping
