# docker/etc/crontabs/seedboxsync
# min   hour    day     month   weekday command
*       *       *       *       *       seedboxsync sync blackhole --ping
*/15    *       *       *       *       seedboxsync sync seedbox --ping
